# NetworkManager/dispatcher.d/netcontrol_global_hooks
U2FsdGVkX18TrSLuiOpjiI8PHSZV5kGUPPlg/JKkKmCIHIMdgLHokgOMQ7k4s/Tr
kqa3781R+nAgjdftS0MjdmS/kp82dW68dL3G+UZ48gZnAv2YxzPKABH0pI3dRmuY
/nwA8LKxsYFnDAzAcZPPFx7FCInvQzJIYcQ96O2oCCt+6UBO3SU/xBI7HAKLSpJw
KQtDTSq4C/bwMrc37UmUN/yWs+xXvcbxT9Bvp4yVeziTsEy58idqH8gCX5pJ0X46
9gqR/3w2mD1A/+TIOqewQi+t53P6PGtq6nCIXrwKj/7kKUPehcX5FzerUzqZdXs+
h/qKzvU4YcasV4DJZYn2kMjPFR6XDpK4UarFSpr4Gkza2q+q9FXC3iP9ItGKvfZY
MD1M1p4oVB/YDphxiHqmV99VBUppV3ExmbfwBGg7qAVxEY7zymkUpOsnj0WJPqS7
zNIbDu6LSZ5oDg0KXEsIksqCrmYwsd8jjSYVY5yC3aMp1LO1mYr3iuj1Ug2j9Gb6
gnDIENoLee8oUO2WkiTKlQ3x5uz28mMIbhPWStiLDsXP6AKoMvm45fwbBem1SaEP
DHfDz2KvlYou5yjPwSxy4VBerkpzr08ayaIUXICobbE0e0XWtzkUWztY3oDgsUbV
O5R+khiBqGxdKi2g7yWflc4txmHeCeEm1qW/Pco9+6vRCItsRfp0F/M1cdDvLkBE
3Ovcenrfzz/MTHzfT0kDV0ZsxLijsWlSMIHfeDyf1L0iD6HZj2eS3AzeRdiHlbq+
0/LHIzpiXm0R9mb9UYqdW+mIElw/te2iy7yILqEhWIHSC4RY/3DQMP2KwbgrVKOx
hbdJUOypn2vhIcsM0Q4+yXl/1YzP/JIYwfcRnN7QoD3OJvzMdO57A1abdclwpbKS
ost1iN/ax8h0g5zwf/4SBJX6Apk8dwT6VRqLpDB8TzXYSE2r+RXg+7bVusVI7CyN
mwr4oDez5zE8YQ98tEQ+4WSU+Sfkb7wdBZNk2jhFidKc7vGhcGms/VtOdUxkb4CX
FmJh+QMmLF8EWGE+obyo5hDpZctAGE7fsFrBU3dtW2E+H2ZxgBr5v/Hv591vW+oy
AVBKPnndu+j5EWInX4CO5SC+wtWr0wMwjwWoXIEVYnM4j/b3c/W4BYGZFHO6jeeS
ZwFZzeoeCVZVK0XAR/mQ2jp7N8oGYwldUOMHK1QleIcd3UQfHn+8GxRy2Z/I7+CM
2Xnyyc3MvRSpNvkNM19bbhFshur4NCj6cjuYzfeid4H/O1UcLKML4jfhzBpdHukl
4Io/4XcmsuWDxj2jTuk2XSZZpStQjuLjaTStb80nh9BNRZa0iFi+NHeMVnaVs41K
5nRKmmg8+cLYSGXLhrdltZujmMb+rvYq2LyBU+qPJO+c1a7iAUq4yKqMPhh5gK87
uiabr020VRAD2CMl4qUL7A9PYP7vzjwzPzNlIQ7ZrhdTaLDRTTrJ8IF+0WldgJ3u
MraK6XF3qkjBgKnSwxNQ/ioqobAcF6w7nFIvNPGV/nAu7doHGD+8eF/U8pga7va0
9SrwC/aUsCmT8OvXvzSM6F/S5v0fWCBLdRaJJoOVYbRhitDYOzPOF9kAdKMmKeTu
W3yFkqadkHyjibOWD78jPV0GNYUvBC1hWOXD3HwiXPCB12naqfeg4YgaGKIPNygW
trrRD8VCwTJybHyjtrom5X0IYhg4zpQNZocz3c3XooqGXboinzYcartkmi1UAsS2
shW99SV1ACCiLSC45h9gCYfyiAAO07pKZuHq3b2z7Dg7e7UznIptu3tJOrw9iG9r
tMl0XqI34d384G6lTscN8L/HPW0Gz+Ctbr1vH4bFUdmg9vyjyP7R27uEG6VZ/MYG
dc47jVeFR16ShENbB1IknfOcLKaVgmpdrVSNyk+ImA7ygaCEoIGRrOMEZZZohpQQ
TGVzHULSZO3n0GbEB3Lplmpm6c/gNP1IazhRYhXCzdJmCoFN+D8Z/EN/aneng+Fe
3UpMAwbQ73BHn5JtQvn6jiTVz03CKebFN/jI8/1MC1Vn5uO7b9qyTvWXrOEN9qf4
6CK0NyYCYcjb/TQ9lBZaCPPqOhrx2aVCOxZ/FIHtLMPConFMlEurb/EMTgd17ZlT
OeYrziaJQGw5umPzv4D7DHcuSpAPz43yDtqRfNsiMZM=
