# profile.d/alias.tcsh
U2FsdGVkX18TrSLuiOpjiOnfkZ78UviOWIcD3Lf2EmYl/K2DQt0gKgfjsrivilF4
VcpSt+F+W3PU96t9EU0N20tBUnYSc0a5MACy/gEXS06DZcXpVz2fUXqUv9dPAfIb
myx1Y++IYFtgTeqmig3X4kGDWo5GveL4WbSbHSKcUMvLsCZHn67Aa88kHaggwvKZ
5RxesHGb44FnE6MrV3k5MUOs7gvpbHtKKGuEhTJQQpSbxudXS858LPbsWtxYvRLr
LVDrWU/Fg37JWI4OnuH8pOUcXrBxm+OBZxOjK1d5OTE4BbvSKBvGvhYJ4yUWZgHy
tAR6kwcr+0vPbXS2XKiqC52+dZJndwALLIQpHV5GRIcaNCaaHqeb5I3PBHi7mhzs
8bwkysTukmYMiwV3gCp42FJmyPMpmrpkctkzW4gW0PI=
