# profile.d/csh.utf8
U2FsdGVkX18TrSLuiOpjiIMuH9eNvHxU2WLakxPnEk6f92uIepFQtMEKURK/q2un
f2NS4tMYDxZizBxPnEZyoh979C113cJJr5+0MI0MnPPpF4JnnyecAcBia05CXfHE
maWew+aIgWcPct02E3r0k7hOSLDSlIfpyzLaEFWwuraTH4tNT1GgdZcK3ph10eZF
qxy2MGjkCHPOCwIQ5rrzESEIGQdldg1x4AjkjiuRmePqn4K1yAdjWMo34m+KDp1+
4SX8R+TV1AlX7/DDleR3PnK2tiZ83cijMh2BYCtsQxbOu6UohrLXjz5pp6QErZKy
4+LipFH2zMYMZ31aELZrDPVI+OAduSVuadLiEWcAfmV0xOnS9KZw4M6DzyROcL4v
KNG6QoNwIrdTasJMq4DTzAhmZN8qXOa7nP5EVLlfuchlBQc4eCCVNTfEYh+vPLEB
XTlPGhr9eUPrNTgQqXATXg9jTZt7w84KG5hy9kLcU7LceIF4IOty+RaBWBk/M7P7
QNo2MnaOKQQh1PTquv/UrJX+xvRmg8sHdEmsXmfRuYv6DtmOs/NMcAMbzxdigDnf
CoGirjjYCqrp5POAMQ80mRJNFTYxiCxVQH+p/7j3GEs4rrvJT/gM5YOO8vk8mhkO
V9Jla6RFFpCLQjZsXr7JqxlPncQfbBkAI1xWGzntCHxd3ot3IoNwIG52rhv8sK6n
gkeLmMnYyxnPEbZp5mvgEvErOrLqis+skeXSR5GJheLs8xw40Lm/LvPMESRtmwHu
swL9Ykun5TnkkVAuzV/1X3l13bOqZLUewKgjMvh3aBW7EVAcsPMzdP2BupyE9Qtz
lgAE8u/O5maFvxCAd1fc2Yl1eZzl/3N/yDeWSiDNH9K3RSYWn+Llj/gGFlaF5Vf/
nSsLOKvyBBTFwfGo5/mugTVg7yYY7qs4CrWOPzYV/lDaQGeSsQSoBS18lObhCdsq
U/lIsKrIA2xXz/xNQ/b3UYEQYs9ruUxzEXjWXfZUKiT8SNm3r9G06BXN527U/gBn
+hBidL7zbPboF221qt0ZOc1GxPNhlBh43p/Jc/KKEVmsjdzxHsZrm7j/ULokAXlB
aPPsgNc+kPZDKQDx0qRSf9HwNsFDdEHL3Qpb2cyAa9mM3uu5REo/vIFBW6mhlCY7
aqAFZALljfkCUwrSMwQVfA01c58maH6DlnSLruzasskh8/8puUvKksjgoDsgpYVP
vByX+l5vZZpeKKPngyY5BTUCBYjmF1mOBt1GV6R2C+EZT53EH2wZACNcVhs57Qh8
PQS1CoNtFHnNkZUJphEEFfugalujTRkDBaEC6JKgMeMzR5nVlY6rFpeuagX/DPPz
WD4qfK+rgGdSWEZplqqK/w01c58maH6DlnSLruzasskh8/8puUvKksjgoDsgpYVP
vByX+l5vZZpeKKPngyY5BUVZ8tdwCtcPEKjclgTC4L75VR1u3SaIZKw11T2G3tV7
YbamZyNk3ShkZpofipLVlGYUbLzKSxDFvo3I3RgcKClNbCYTlnYm1bFi/Rcm9/WD
RkajLqqdOZ3wmEvFL+dE1pPpqo6N/QVr5hPcwyv+k49I8xWXG/uEMrf3rjI8C0PP
