# profile.d/desktop-data.csh
U2FsdGVkX18TrSLuiOpjiK5GvSu96CSSSEcwaP9Ebj8W8l34eHhaxIALa+RYjsXr
MM7/XatsAY7rU7d8hZ8/umMDrQMztdgR9r0oNvZheaXuR25yOohN97luCPQzIuw0
c2Yt1vdFFqetn1S2ioYz+ciGtlglc9g5mRmzklEh8WgMHA+XU//Gqx250swUgciE
oeujoOrtvfXpu4I0gERCZmyvgS2JuzfLNYEd+kig4vKQ4BAgQTCKtH84tUbk/kNV
lR9AWN8H7959MJgzuB9bETR5wRhIYnKZ/6ezZB+RXGS8HjaRCHqQeZG6p0IhwfeP
6+BrN3FYMKZA7sx2YgZ864Cl3EG61eOhmVKHGWWmDBq2gc0/9X1kwu9YPHGZGrUs
hK6zUCBwcp8jwAPV6poBexFi8U2FYu38/euDCCtmaKAJK/Ldu+lmkb7rYYh4uNYU
AnxNtIWd3XTfqXQBKquGclwftu8CnD06Wsm8KjdWT+qr8dTZVAC2FaWOKYqNcnTz
UvuzsnR6hM5/UY0z3S4gsk1Vea+tLKBWwnedexi+dxrC2BK5GBOqBaNkCQ3amd0D
YJFsAVdy/m2bpYurqbtEUG2EHjdS8IHvJkTgcrD0nrF6fHwNcKVETwKuNmA2peqB
p3Fx4BQJx1TF3X8oyNwZZJj3nsbfzjhAD8z+Mi3tYmXro6VPMLJRu2M6ywpVbhL3
MtbOQj49dGPvYFzCBCdZCS397u8aVQdjTenl6qPnRExEdqDFKbMjWlSBRGDbReIY
4nDcQhV37JX6Vmh9T/hkq9Yr/hFKO9xcZCO0HRolMqwMfd4olGwy4aYdl7i4Vffg
Qbu5QxtLij/VANja9BButFxwo2GqqSDAJe3qv3+2G93RVN2kuOHCZQ0shEFJ6VBS
Arvc7+7v6ITjhjM7hsxtKQ7lBF3BKUsE0ywHOMs0FKiEt89BHknDl/LAbRb7pyWq
NswPEn5R6DtP/xaBc31JmEJxpxvvjnkQmU/aIGyp27J8KgPQ7E2oKQmlJ4fxd6Gw
uIsISpWxoA0HHFaUcjByPswUbywtwx6EWIBUfPnKmI36nzYdJmM5SiIJn9RgJge2
D9LUQnU/1dSnzYNrRO0rHbKDurKruKcbXdw94RQYyUqOGbYI3vQ2rWX+mNP4RUjs
M08qhGQMpuNxSr6ieNTQyyNEtT1N/TISruEscbrSUhTBDJ+9JxMpUhgbTCl34bIk
pblpwcdC/acw6+EBZe0ZkmN28oeQlgWZTRuZ6Y6+cN6JrSOEBfsZDdLtCKQERj3+
vBzAWgIkhYg8l47zXifkQg==
