# profile.d/desktop-data.sh
U2FsdGVkX18TrSLuiOpjiK5GvSu96CSSSEcwaP9Ebj8W8l34eHhaxIALa+RYjsXr
MM7/XatsAY7rU7d8hZ8/umMDrQMztdgR9r0oNvZheaXuR25yOohN97luCPQzIuw0
c2Yt1vdFFqetn1S2ioYz+ciGtlglc9g5mRmzklEh8WgMHA+XU//Gqx250swUgciE
oeujoOrtvfXpu4I0gERCZmyvgS2JuzfLNYEd+kig4vKf/SrQ7Te4l8Fc/Lvk4hNB
gMoGnoZp62+DdKcdkJrTp9HiM4lA3j6e2XPjm3oRc/3UD6jRGwP0BceRaiDri+3Y
ipgvnIKgVjiWiZZboBVUaAuY2O3GLL2NYRp5dPg1EmIYwU/9jbZZ+4ee/zw+Dv2V
gptqBFnWlg677gjAXYEPHD5k43JVqYF4jOulUWRKSRT7X+hE9Ridn9QCLFuEfbca
iI3ftMHS38bHpDINJN9Jey+ph9B7PoLooRGYdRKPlmUNrZgUcaKGIvJLPjgfvJj8
/V1fSYymbZfFgZjqwFPuB/EHQNK5PBLxTTlDNxG8DL9Xx1/S5TiWqmoDroO8r3y1
4Q9bDXNKtMjax/MqqPHLJcALXzzruaiJM95qNe4BoCpTQxKYwnyUkDfZCrkCnLxw
ol3FiL04aRKsP88ogafQ39s5gfi4dnK30J9bR475L27uU+gEjvrvBO9YUz0b22hs
cYLZ0u/+wl04WC+CDxfxWbpuDQC27621sYKiWFzZMgiGD4nh5AsUtacY9Cr4o8b9
VDFZe4UZI2y1VaDLRJee1+u0iRLDNt87at1dzUCDp6MmR87CQ7kOkCOmXrC+Bq3O
Py6goMvs3pYaVW9fvk0Zye8pUl1UweDExVwvCACi2PI=
