# profile.d/lang.csh
U2FsdGVkX18TrSLuiOpjiFd7JXzqrUG24pT7E/3meyiZU+MCvK2CBeAmISxAWhM4
UZa4/VxA8RsWI0yfVewFCBJ1vtpH+fXgS/XAWQtnKO0BG5Q2S1uL+aXSFGE3T5wa
Qh0LMM+G6EoV/dx7Du29JKwL2erO96sNH106zUGJxfq+9ov1Z+PqcdzoWIfbB0LS
8QuCjCjIO8ECjgEhSLATlm9jgPqoSr3RYYSKQTuxP+vJC6M2D42XY6JRZvjb3Xnx
FG02cj1jwo2xxgzVcGu+LPg+TYsJbRWcaJZY4xCA8F5ZyVe2MGo0kDKT1/DbLJvI
Q58ntvvHxFACvWDnBz1Rys25qIbi3elEGV7PDegyuaCrV8/ZdfJuY2w/3vjCXEOB
c9VmZSMr8x/ezaM3zQS4R37hKJi/RK+fy9xKynJaTHD6ZkVbaNqaiLwbIDQog0Y6
i8Qm01eTYVVJZwxx75JSD7R3Thy1qPYWjp4ClZQNSyM9XIdj+RH+c9DPZEUOuW60
0n7dd5ysilJwCbIV839aIL6lManebinNJGpYDcCYT25+ELQRS/isbNnrp3T5mvih
1xsSqrvaeH7A+IKGgkpf3iMtMTGvArqWJ2QFYysnPmmwqzDymiA7eNqLVJykDtbg
0Iir1nQry8P8JOj4dKfyCJoXYBwvtTMATH9mqWOXaebWXwb2Q0d4/sKPzj0sHv3h
gVAiMF3y6usBHWkJZ1u7Y13Jjxlt9jXiNaiSFM/rHSWtokwp6U0Y/2cuhuPqLb+Y
bsJB/mJmmdTpTEyNF/vS50y4jf0L72YqQ1wMkC2YNsNTtb2Z/EkyNVogZZjtIlUH
3JInuwJn2e4gDwPgGPoQDep2HFLzJkmnhq7ZVeWrnOtO0PeRM6RaeZBK2HuUm2Dp
yHf+TsEK8pI8ZEE3dQfA3MdmOaqfseLNrPTKr7p+rWSrE5aLH7foNwiBD+ClUmy6
Lxd6YnBfS03UAwnmN2ifPjY6xWv3Oo9I9JDqDvqznaZbm7G2rHWMM8c78MSf6u02
fGBV83mmp7/kqzpI5ViyUl30qCz92ezEDNzoaXE9MVO3btIoQpIqRzziYr4RgEFD
3cp/lz5ozcTr2/ZrYTGEtOUNUGQqdykYSRm4e07NFWF1yzpFHOn0uUZysp7tp0P1
9HiCG3+ENIIYWbEUoQdzD10caukDyWKadr6nGGKqXpLE36hjiyEOr4xWCFo1jcep
C9gQxmcgfkcxP2fVAs/5uRN5IDjN4pfcFcgvNHv+6xUmB4SnLFYVGgHTgcmgz2HZ
6AA5NaUVPQ8B5I2ZdrwrQlSOL2RhI48SNxQmzTKQbWiAOXrNdW73DTdAoWsmR3qr
AT25sHIbh0aNYY/ylfwtCYfqRST7qeatwobLWUlLsfa4sgi3xDUPzoxgPuxtzKo/
QdQC1VN8yPOLcnVL42dkKmehbCy9p77TYOmnpdLi0lHYKtULK2sDz/FblPLacbyj
Y6F4MPhyCvpB5n4ZpVikzn7sOq8H2zAfCDJx3xWGZzOIWmnIZEFCsAddTPdXXDPO
Tmkd+oC6O2n+dKLvhImnra+qv0E9Wi1Xxgtiz5CdUmhCWpAMLHIYrDGNmhDoPpiy
cIXTTVET2ZTqBO4TUW4pd1CDuS07a4oG8Iz3JTN7gQo1eyb0nkjVY0itBrjAOrNs
vJXXcC0WAybLIZLjt0D0XifvilOVLOXA8cSSKj2zVjQMej3GOprQDn5ICTFkZUxX
X9gboWhWsP0CGFRduxsxmU1hQvZT98Mu6pbmJRDkA8OfY3G1yC2RjtxQB9y2Su1W
/O8hiHKJAPkiQGTyYORGf+cennDf6ScVqHG/1tRQvd2UiCgG8YB/TEIBtYWf92bo
AtV4VBowWMaeZUFOzRrzNWKbpNChVBmYlzgY03MkYkipPeucVufaH2FtR9DIB5S4
nmdb0SdnwSWX19f2V966mCiDHbpVh9saVh0CnJQOtQPPACXLephHCVWF+LnlNyUo
Cht2vF1P3Eg6nXEkCVdDa4UpV3esKe2Mrk8IsQ/H3RdYGOJKmul3Uy5FDCTlK20w
bd1r696SxP0zyuKYltM86y8XemJwX0tN1AMJ5jdonz5raPyqTQcZ7di5s/UxgK/F
/0KWzmAZ5eipn3Mql1q5HXkzv4rlLW5bQaKvIRzFmSJdMmgzc30/RELNc0c2HBGh
WsCGlTPKGzjdbL5qiI4K81kYHc4uEEJwvaMzV7bs471qh7eyWYqMkuNJHmICfoxP
V97B5VdRDPfM2RqQZbzPnGP9Jt8RzYmnkumml7rOVBAdXHcOAjoqCrS4ISF7LI7X
BkCVMScx5RjMPg5sz454jv4bAtpI3PGlK4cj1aEZJYQM+d+0W5faT7wY/bEyG1nW
pHagTwN50JtfVow5ccw6hUhPtTUmqRY6dTqOTGvZ9JI=
