# profile.d/lang.sh
U2FsdGVkX18TrSLuiOpjiFsByvhdY8UeOoQNF49FSvgfrddsKBf/Pqpi1N8G69zE
64cQ5iGbHzbYpN0tZi2JlKXQjOjShB6GfTGLO6k/FeyNxMwoFGydG+qOTWZidKgR
Uc5xNE0IpnjXj1vBgjEAP7Uuz1l4QpsljdBugaI80YE7QJpqU9pPFjCaYk3AfnHU
kM1vdDA+oMngrpf8kZRzQQHics304OZFYVxTR1PIhib92y1PxZxSmNihh8ufaDv/
82BueQKg/mWqqaHYNHETRJbzmPPM/jOxXHrG8/9hqJsNWziGmw7f3C2w5v4FlwrA
Q58ntvvHxFACvWDnBz1RyqA5dy/p/ald9VecMDKSa7qi9C6kxxbikSjNTjcbgmP7
FMjFQpBf9Yltwai8RfDpDSyI7E8T8ZsT4PO0h6JIltoEj0hB2lWWDeGUzzVxJ8AM
qpTgljN5j7s28SbAxhVBxpKOBdD2r39AIOT+uLOE8RAK2XLnciOIfAca7Y0RdZj/
KUhHD4fx+cdyw7A+Tf4GTv7Aqfl/SWoxaOBfoiXBtiy3TsM2E5nZoPmS0AZXpj1r
Rw8wgSt7y9lStcLsPRFxGeVpZEuhvSVqKzLnRiRNGWjzWYmZTASZorC70N8w6F9T
viBAOfVKiRow48px4CkBnVbXqVj7wqwVYZNZRBo5DxVm9Z2htwbokLe5QffUWN89
pUvwB59jj92zQ9a1384Yd/Yswz/D+//Wn22X41g9unwm8/N0SCl+LT0haB1W6lxa
PGcnkyteUi+uc7jGoOh4Re8h4GJ2W/CMZNUUTgTqtVXBEG23VwatakhMK681sW/V
8prwcYQ/ooA6uSnXbSUj1Oq41d1mbsJbFnts50oA/v1q066p5wp2BBfJLCzn7NeZ
I0k7e3IVuEBDZ5ziI3tDMcPV5q5Q7STzox/ns5I5OS95GSUpFgjWZ8ZHZDrEjaQg
8iHDm+30JsEMzNAVEQ+UMbTpGWImoEWsiBW6Cm2pZbKFZpgHGklJNbtpUBp5XUpw
Ya1LMXb/KSNJJVpMjE12yQ2i+nljOQ9aZ0UEc7rxLqHlDVBkKncpGEkZuHtOzRVh
xPX2UP3f/EPpN1QjwtFApLYXHed4Lv5UXJP3eX6h1Pa3RSYWn+Llj/gGFlaF5Vf/
t74vaeK/lvbVWRV2NQTZWM0aXQp+9ISGh1lgRf6AbKPaFbD0nI7oeoCC2YZ/eGbG
dYwJqLRAJDUdLjhqGw4FkICV6z9MESam7FcrqN6BIJOIfGVacpOVLd7kRCMWQP9D
vrlfaVoBRB5B6y+MCg5sQcqM88IrTs/Bud7n2gifr5Ga9exZvGdWWaxgl5MR7pZd
rLTPLnXtg7sqhP8cjnkJQSuAV4uvm/wS8bPPooo8G28NdD5/VopX/rhNqjLL9iYJ
guWvZlIIoEmUq9ZNiM0n0jDxKJqa+YOt5A3htMBnZLcMFSO0bJU3EJ90eXOjrvXT
WU9eBMTNOavFbej5hqAyV+rKjuK/mahs8ueX6pp/iho0wukRix/xQgJC0mm2UyqF
CMu0F82M1JLY+SnxCKVsrSSHh3VgWB1kXNVXeghmdH8BGA73r0Ilu9uXWd6Vrsd4
25o5ATrVsXoFMrSvO2wmCk9A2LaryUZ7o2Rniipnt38WF4qk5SWgrLa6pv0ub+Ok
1vBr0Jdrt4C77oknQLScHvoZWUCm556LHGTyAyRnO4jQcOsHzINT4mKatSbBJhuD
M6xgpW2k/qbK807OPMYMJGFkKHCr8uM8sdkcE/4JY6FPQNi2q8lGe6NkZ4oqZ7d/
b2vfmyW4/amLzeV7zlfKAu2ulEU3M64NBjh+wxE9uikzwDYOHv7k94430rDRu17q
Xtt9UF9kZ8e6KUsMvdFiCfxXIuwHB1l35CuyAdkcpH8WKGVmdeHgaAT+IYJJNWAN
jlXLPXlImh5jayVhcKGZcDb4KCONIFdbKLzsVRZA7nvptsft3z+ikJGQBG6wqBYe
OCePpWjiD4YhhNB7wZ8L9Lo9bcIkMMDfkjFJJth0wcbaJku6Ct35nh6Ai/AZBiwt
vxcW1O7RE/HQJz/x5z3vr5T8FedeL/l+t0+Nfr9p8t22v44JogXbtqyUJEIWWlkM
oorzxUADvYZ9CGySmPw7+r/YMa152H/4aN7E4/ScgaxMnXO2PPOil8NnGKcojkPW
SPMVlxv7hDK3964yPAtDzw==
