# profile.d/ls.tcsh
U2FsdGVkX18TrSLuiOpjiOuZKh/CPtMrokhj6hXlKhZNTdnVwUZs3Rz6aCFtqYJs
Q8OROyQewli85Est4QTUWIswAPS1qDyGeUUQTGKzM74pA8WCPuHYvfdptTrb4kiJ
S2cWWLDjwnpJu2eacsJehKI0ejP9C75k67WGIJTR9oJfjqrKa/69NGFRx+oB3Ixt
h5vg0BDDLAxvQn9oHswDTihP9dZRFKdTzY60hBhYiueo5SYxOsowaWviDmuShj/g
/Ychyr54xAdXVFyrC0KQCalfTyCr3Qr3MOdsllj+d7GcfXONkz7E76i6gi2wT6S0
WO9IEN6ayGNkoC99uQTFX6cjQhtYIDoXmGhnQi6moNmw/xCb3lv4AQq5fZIJ12Zj
XObWN/ASiACtegDEMLx3QlPvfqcCbKUrLsRMgBTqtJJad6ZsK4V/RieutIcIVWXr
8B02HLimvWCNloSb5AmZMslHXcmG2n3p3CKQso3mucy4ZtjcJzP1szI5wBg6p1YA
C7wIH7lqG8f3g4BfxOU25MNjpRzD1mp5RgBa9sG62lYuTQjoMQyC+XDi5sV/VyDY
qLVmYbMdS2ydVwgbH+P2+WkxyNjIRQ8Bvq1pBoQLPa7IgK1YFSvSB7b53nNmK5tz
V0b6ac2DRbbDCsdCa9g1k6AZc6Y8xkbSIU1uLlCuGz0SC9PsCZ11mBg4NAmD5oQS
RttsK+WW+7WpwQ6CegxbEC/ze11+JB+JPrhlxGcAHGNcUMXSPXc1hv3e//zkrWgc
mgzYgT8sMAxW9YmSGBI/eEXBTUnRJCi8u4iNtlzvRnDbcNJl0rJpEjHqD+eBqhc8
O8ce6t6mnh+1SZCUxD58LDCdCJuzDXg8l4LhZnlP+fwCT/eNyytM3EqXTGYkb1c2
g9OwvE8VwPYcy85RntJworbApGSGle8gWlFSwINepNvMSV39nDCYqaszybT5xlUS
HUZSn5N64P/XTCmOzBAo5UjzFZcb+4Qyt/euMjwLQ88=
