# profile.d/PackageKit.sh
U2FsdGVkX18TrSLuiOpjiEdqmP6av/z719Jn6wY9DCOUutzSu3sTpakLpPQhhOoN
piso4eS4Y06gtZwT2GTz/UhSCS7mX9RedxCOYNBZHhHsm/jHqD0kb+4GurOuMdMc
Ig+h2Y9nktwM3kXYh5W6vuDZatd04Z9LIAnAHaZxAi4ff3yWO7H0oqsYJx6kseYo
VMm96yvfbuPwiZJr7ryH2MhuP4GyUYHXK4SLySwoEL9PE2btMK1AkXfLbASzV3T8
xnICncH7VitAFpTbJqtebYVnY+XVErOsPS6GSvnnbvczrQhy//DtO4qF40rllOJ8
uy1eVzwVWbLf5N4ouwuK+g+wdGgo4FwjvSYGqs8NL2o04vqpWE7vE9oLzt68lHBB
hWAKz0tTxM5CGS3SjNSbr+d1ZmpW4hAL4QAhZrUtVmTQ8/8hjRkG5BQIISs4sDnY
hP7jEq0rUpJFp+UttxBJRUZqU+DtPDForL5NJMY/Frol1WYP3mQF8uPFCSKgXSXR
hDhgcHdQJTl+dn/bY9kB92HZpHjGdpae84DBTvMWRqBR47/3fPEjKU8woWDszbue
1vnF2QX4LRWbemejbwWpxYeQ3Or1knfKSpYskjV+HGSiyzYtfw6aohtSffku7kji
M/5rPtQ/wEELYySHUoUQvlqlj/e/IxBv1H6a+I4JAmeoHqQp1Ix6kwMBx3bhVz0Q
dlbv02L0QEnyxtOHIJ/E2tiUTc8S7Z81sWsLoupiADLbuME0DCrJY/I3fucZ/zb5
KTA2QwsUybuBu3Z2g6XA8iB6JpjCuLhZ4jOwZzCnnPK9DT+qdqBK0mMK3NcxIWIr
xvZEl0nUcIg3hQ6hVNshkDU6LHiiTEtKbwJR5C/sBC+58+hJp6/ojfqoGVfaSuA2
kBF+GrATY+2kknbReDH6vzzOKCs7flVg37poDuUr5u/FqCKh85w/ZR0xqfeE8a2P
8bcSSt20jwBHh4uFjsZD3PSzBH+rep9MCst0bsY0tY53/LF8nxtIc6G24WxnoY7K
O0NxXrF6XyPoP6ZUqmbRV/msUM/uMcDTt/l3JDjOVfyzJ3i080GGU2GHptt5r9rk
5MUKnyytCbiZInB/Q+wq99u71wvmVtoR631/HxHf/IzT21MbzZFqQMeIWvpwRrsm
