# profile.d/profile.csh
U2FsdGVkX18TrSLuiOpjiBB/rJgaY/lu+oJ0UZAKNjmpRbdh6Vthp8u263YjogHO
cpB06ejTrcRsLOP9ppygb5Ux5v+jr2F1Q/s0N0ZFq3NmNWFvr/jDDp9ktimdprgP
VFClg9I59+4ZfWHjdbEc0IDKBp6Gaetvg3SnHZCa06coFqEnzgwslcpSCi3G+RPD
wD3/vOYE/IVRkeqF1/gfglQhJWqnYkfh0ndG+Jyn+La06FuiWszmUTaWYb95/xlp
VKsoKQ4BtYp8uHBLh3AxS86I19tawWAp0nNd0jB56vqji5a1zn2dAyCcnPTixq3d
dmAKhO+QQiDF/oBR6qzqU50zG6kjvb2PpnV+NB2XqEY3LxytwKJDihz8V938iwCi
NB+mkNlxC2pqprBwdeK7mSebR5x/wD27mbiyLrrYtls3jWBaPR4ld00y9VH0gBgT
Y2imQ7VoHpEUUYUxmMb1o6OLlrXOfZ0DIJyc9OLGrd328BaV03EHYpgRAd3Lgfsw
Nl1uyyHnqBdRwK6TJ5pCOwS/THADc4VMKToTtVWiTPGji5a1zn2dAyCcnPTixq3d
nw4JE3b3Xnig2yN1jVK3IalOnl5t9mtOt9W6bSvzjNCmsyt2XMNm8SjOQJqolWjl
qU6eXm32a0631bptK/OM0CNNnlJupzvLdMNAYszopmGHTmuM0OiqZYJ7DNp9p676
XyQiMY70+/CdCDbzpG9KLrPvjtc+CcAYcI/RJsWZ5toCOBD9XHGNxzTzX/WpBDVU
vmdPxfgKix6Qn2w7kwXzIvXxZ7zQojXYsbmC0NM9OcTdJvm2ct0HM5r8bHvVEYSf
rM2rbojzeUJpEqh92LussTiHI+8sDYzkKXcLoO84dj7yxXNZzuzYz8SLUVgsB6Xk
iJAlkGNvqfQ+Yj7RsgBd9WXmzlAVB/SZC+QajHqdPNosnKHvKI1a10uP9GUSZSsq
6XlaBwoGtz/5/WWgS3+odR4aGMPZMXhslyUtdenke0U583461Otprh5KFkjCZ6VR
SmAnzYu18sviYavKE59DlA6GKQX/xMhxOv7bdDsWYLox3WMDDZbxvBuZAidx+2wY
CAY6Kdu2FcId/Su/ZPuv7YQYQKNMTmakEhOpXEeCnCNQg7ktO2uKBvCM9yUze4EK
ksvv/nDsSbuYe1JUXMUy/ko0qoDp7cDG9o9N3ANNHrwZYZFRXI/mo7fA7GOGHWIw
5W81WBjZVOfisHAHN3/FsmwIdX1I5/+69nHuQKF8gW8tsET7MB8Is6mRTbj/zrcM
PqyCyVxM4f/dZkxyh8k0NG/Of14HwC3zNiq8+oKPyk6mYJao8YnVbVsvVNr2OtLv
7ooTXE8IdVOikLJntlOQIbPxLV6EQY9I0lslUq01iAxBHLBlDXNEc8gdTCsq1qlV
P2cAxDjaXIWeOyAQt6clkNeiVvFwFlOA91ppNmuSQf77p/+HTTnNtWgzTEp7SKae
kroF/NgOXIvct/s+xWQN5lrJm8HaOHL2Q4Ho71E7959yZ4n2PCBtPgbLk2YyKO9j
9GGE8r5DB0xA+X4xpDuZZstvI0SbZZp6q4slN4x9XYFayZvB2jhy9kOB6O9RO/ef
gb5dps8X4q5mafPCyUUYO0UhR8GX3r2SxTR0r89iNK08haXryNLOhdlHo15agjFZ
bbJM0zDJGXMvWbX6EwQ8Y0IY+lj169WpzMfPZEp+oC7DbIei7MeN1gd0LT/fCQDZ
tpYlIdtWuwN1aN+nUmoiXDQUHgbct6wLlwrOm37f4ZdNnPydCcWdnDlweMlBvKEO
9WmOyNG67xQ018FLugeKXajL9bcFCSz00pL0RNLAaWAX5Ge8jneWl0Y6k6we96AH
xVKJWVcQIyn0wKKUUjbQRFUU6sB/dVsiU5Rj0q6ZVqzGLHwKrSKy0xBy+wr2kvRG
RpgpTOBr5wuKbLhqKaXy/St1XfLwGT33Ut8BUMxr3IXaFU5PtphdK2XOfp7fpZ3a
5aUmbMPWfn0wfLdsEGiBSgp7LCTK4lDMNBJ4NvcdhHnpqn7RivcC+xlgcq0NR7Y8
9xlBCSAyy/o7PI/c5oHXM+eWLJkAPpyl8FxEAZIeGLsGEbqf3woMIsNKdGEKZjhk
8gGzoQ+uBakCJtsXowYEYU91pAdgCuGr6jFkO+5FLpOhQ4luSBemhNiHAA/+Rjzq
VH1x55kPWU+D6jLRBEzSM7Mso+w8CwWp+bUv4k7Rrp5VFOrAf3VbIlOUY9KumVas
xix8Cq0istMQcvsK9pL0Ri67wwwsxSdjV5fX8mEe/RW2Smo+j2oaR7wSLBL0ZkkD
QSk3Au0GvfH2/ghWUzTYBdfmj+//hZbDW3vCCuxXfzvAWeQYMQvNjaoWdNa1HiOp
PdS5Jd7ikXhGQUYTaoDlG02c/J0JxZ2cOXB4yUG8oQ5GT1QFapfuDYvUmI4qMwQ1
9UrR4hv8hsaxs2fgtEINyjC0DHepJ+DFJu5nLlH7s0tQTTNI4cfboIbxl+YTFymh
TROw6COELUmAKEF99CnF5G2yTNMwyRlzL1m1+hMEPGPtNnCGMGULTe7sUQeQT65A
zGcCxKCaYIrGehz52yMyqsO4f3qrU4ARPWVaVd26vZBfH0sJtCgY9CsUNbATHxGI
dymRk9HwqAzxr0gzuTjKqPEzWikjQpOK4eW3rdCPFxXC1zu2uZgBdJ2nU9aTVnaG
LATsSo8pTa3aiFGy0wYzayDhrJLR8H/ga5yr3kekAf6qm4WLko157DCKwwhtajk3
m/sa575CVfgcS2LoAahDYKPMqVIa7PFhsJmHyu2/A6/POtnnzFNbNz7IA7o1wzFz
RJjJ7b3o/oPF7o5EJV01ZIieXsmF/yiL6x1Va1CPNzGuV3egK5Y7/f39a2riW+wl
ZnWUll+v4zIghRQaFO/x+F83eidzbanxWOOnM5d8JfRBTEnEUKEa4bOL9jyL8qUi
NePzcivzzjQaTQen1qU5wvZf5+XdACpZqy8HxXWZqHHPdO/Bazvq8U9NIvk8g1PR
T9HspgUbNEiot1SG677CZAPLUbqpeOYvFnbeJxEUSqJEZIlAyBAVF1/g/wJCAebi
I3TSvLLzQRtYPCaMRq2zViiUFWTtTMevWzz9Q/0ps4exRWmT7hUcLMkDxcPJOMDx
VUTrxompfx8P3ef6FVgq80oz6ejE6bS89GBPmPULQWiKsjdAI55n2KSLm2cLgYM6
KhvrMop5tEkXX6Gw+MgZLBi73V5CeQv9jBvfiQxFNSijrUkcOYuWPBuw283bt0xn
bvLZNKhOdgHb0QSV7o1UnULrD1pE7e2771IoNa5oaLvF13ebEy5ebvgjEIQWGe5l
BpQAxIWWZyNOtyDnqwlt9Qxb+ofMMxoRb+vbyMcGfHxJyEqLfPO7vmbN3lEJkuj2
12ksnHWZkZLd/NuiqGIdZyarDYVUFpKllIauZ6oLDnraPrHdSvGKq7fZnyli69gp
qOhOdAZqk/J1f26TXlB98LZKaj6PahpHvBIsEvRmSQNkfFkiVX5x5fVA50QVrBMJ
UXbKPcU8ptO2w/0k0Yw3xinnHtr5ebARVSo2dGQSvFhCyFusEqODrSUMyKXgSAvv
RtUK8kSXMIi+m8QjYCU+HL+GNAPd7JvrveD9QOcAcjgrQtSkCpkNNonxqbvP/P0I
HYbWE8QbparnYYOEphe+wsj6bs+CPf1Llnl/JxjT59VfV2Q8spu2KIMR7SCgWXZ8
FBC2Xmn19FhCG361JsYo+89/vSpYsyA5WF+Murj6LDcbfS55I9ujvvcdMUeCAS9F
6bgQ50ZP/KUBl/yDIGZNv1cQiYdcDjb2EeVRNBwyOCkim2yBYvYXFMo1bKCW4lSO
+okHeovJSuybuZXmHO8XMNb3/9RCmV4W5+0/ZlVzYgOt7385MLjULzr6rr0zy+H6
hapsyGqofbqHqmZud6YbADdOazaUh7pLMstkkLmOxE3sBN3D7mfg1Sz75Zl61ECQ
ZWpNOkNQlxNlZOJsdBJrumEuWH0sew3aDch6dyl3MsNBUfnbkUXd4vvFjoV6t5yH
KhyS+GDWIZ11vb6oZcH3rJ7x89Fg9LOHYiWlw6VFX0Y4Z511PREDnPFgKhZfv+2f
68n2AvrwOmMn8rOw0uQxvoVvUjiX5AV9RCVHC6ENM0gYfEkdzVsjaDERZKCf3Tt4
YMeSr4u84W/gtj7h3Knh12Y49J6gyigGHEdKa0kYDff8pECPhh6hhyxI/jmcK1oP
JnD8kDPpIN+2DzKZ1ArwUGqU7OFzMLr0vV8yPc0leSlCAKWo07hbSQ4UtIlB25fM
kpmzwADplGhWRhs5g3cWRdUIZb9UhaobCQf70NoSR/rYVIPsNEl4VOtU0PrMnWnJ
k5tl7ybWW5zoshGzBio+BGNdg57eiCgP6j8osi9U7wF+OYc3MZxYOdGMB4738WRW
g4wdepEoDB7uvrF4rGWydmEkr670iExlrjqIQecWUzBW+hinxDHVp9ZZC+UQs1QT
TAx42yyy+7AHrsqTDz3entnlJAinnUnk4aeT1SmpldEVlINJxVskGVz7Y3jmTz2/
CuraLikO8aVyc278QqvOU4UN5usfrHumldn1+Od5mMJlQE/3BTD7hcXYAG5yhpGr
DVDjVASexB+PdpgZGAnIPb7qNU/oOWJlpBjoS91L/ad8rkDTUsz08aKbRUbTMOJv
