# profile.d/python.csh
U2FsdGVkX18TrSLuiOpjiClvRPWltbbFsyltWjGNFe1iksBpXjOYah779hwQyk4r
8feQ6xTCUErg6jc8ifkn+QRrkrFOqf2bNO7COuwfeo0JWw8vtlsZ/F19olPEAKLv
TSzH/8kf0MWZgB44ZcukVA==
