# profile.d/udisks-bash-completion.sh
U2FsdGVkX18TrSLuiOpjiMro1eOOFKsyDpwRhCSt9OYi/S0OSbrMEvCwbihcvFqi
pCFV6AZeALLtCuJoM/LsssJvMVJa1rL672+Z45NRFffWOWHUuSjlFQjwgks8FSQa
1jlh1Lko5RUI8IJLPBUkGtY5YdS5KOUVCPCCSzwVJBrWOWHUuSjlFQjwgks8FSQa
1jlh1Lko5RUI8IJLPBUkGpvUGQ7GWlkANxm7LqeiVgpHRx7eGd/KyhX3u0tKT93G
HACulWY+lC7QsKrI722kRXhnqaYeVK+IWm0xPBmdyHRgu6IfLEHaOYUlg4Kr5JJW
NzdzSkOV72OpBIlU0Nak6o3VkrN84arSbsAjGKlsRPfT+HfG5R/QlqOGZ1ZNnMUD
XP3Yr1fUBhiwOtQlwVGsC7wPj2fWaJkEcrLCdOb6IH7HxEtrIHqUWAxTvsxNj1Oe
2b9Ns9pxpDspI/zRoaVc8y3bzYuBjiEihx00GBH8oYad3/DL0r4Mh9rXTucIj6Cm
C0/jOjEOOxaUUiySuYmSn7qocCqZY/Le7mwx/eEF1/jcfsZNurN0pe9HyKkOZnfa
81AIKWHVxBtUbyqM7hlQMExicARdFuzPtf63vyp07yTA6HAQ8Twgjoyl7iy75AzG
enWQJObDu5VSNy7G9jACZ5rsU9aOVZdryowZnJeEeX0KXTOHI5JPConmh+U59+nW
XeR3Hkcwtju5QqWMFfx8p61Dfpg0F8Dd0YfhxXZAFr9ABkNWm8xWcCN7HOdnF6z9
QeSTRQNjCT5bzDIaNDOoWvOPHuZYnBrpW3f6J/+3StugvwrB9h8kynmlWKpJAENy
7DKd0vanpVdjnTlfTIOYgNG2maiud8x1yxpUUsk9Oz/oHJgk183KBNcFxsO4rj6S
kI3ljpqunr1RwYcsr8Kv3IinUQiXzTVx8fr8RHq8a4FXzoU4cgI3thrDB09/mXuF
/NuQyF9pVg8hzFvkY5E1R1nc7VJgOc5fR2liQEm9kGMprH4AosjrpQ285SHap/X6
S8ScuNeQTVchm8mqI8reL5YvJbcVB+YphbQHJFi3GHYAKnMYzeCPNDesWyZWnBZ+
Q7u1omNakvgwEd393lon00meKSoMNpD29laJWr+IpDNOf84tw+GqLEG2MTe6bL2z
Ef1INSo+NlcMLW3Zs+DOcq+F0iaXwIQWjs+5Agt41ExY42MPey0oAJTejlrmQGA2
Ag5gDre/XNST8sHuywx2VW30wm/LGlSd3dsR8j316wxgofYlVKu3qr9VAGUPoGS3
9gwEgdFLPJjkHJ43r1SO/FMSjGdEVBQiWJGWB+kdX4gevxmwg7JPw9pDy9Phsmez
egqOaSR9mKxgqBuPAoy6/mcqOu4sfW09A50Iqy7AMy3Z/FbtQ836jmu9Is2esx4O
dC1aLKJWONGW6jAclrWZKP5y0t1ENtHzYsj5J6ko703l+aOHnomEUP1wud7/C9jK
EaBFHrcuszDyoHxb/wNdCaDwBCAkxhpaVSHEZ7WN/Vs19PnENKTAsoYLPaWStaTf
2gJl43HMf9QP8nnQX5mTZBesjacUNh9scngy7ZbfnVR4CccD9ihJvnMWRlbW4DNA
yPPouhJFdCu/3ZxrqJGYmL9RnNnji2HgGcIUsmZABpcCgkAX9eHOEbi/h/GKy0CE
lN4cxfCUyK2AhQaXHbwK5vSfdK4kIqON2GdXd6yJNdtTZZdqaOHQAbYZ8O4P4IQe
xaorSGroMbyBhX/jpgPw4dKBsJ28k95YyVhksceo4cKN1ZKzfOGq0m7AIxipbET3
0/h3xuUf0JajhmdWTZzFA24+OlqYJu33jtNUy4+CPE1aFwwMSsSgaIANZ4UZvwyl
q/WgzyqafYGxFVE5Ng2x4LtqFe2rPG9RHE1q34VxhF5HPXjN+o1Ou1whBZCYNNNw
jdWSs3zhqtJuwCMYqWxE99P4d8blH9CWo4ZnVk2cxQMRywR7nLcVPyNeLh8eCUFr
820TjuqdS6MtsiAyqai3Pr9H6M6seT2VzJY3oCDiVIzegq0Dlo3c3I+pQ3SyXcGy
TnJDfDLwqxSHIZpW6w29vGxW2P+8TL6wsF9ELzEYOzZIORZD1mRXWpm34gXLK1tD
h2SaIqG1rN4tCaPgUUeIOEC5DoFZf0+UIqa0rEp9Wgm4eon7RkOj5Gv6nDW8yOSw
5B5ogLAm3frxuxrzuxY2p37Qsahk0ed1QcXWlEKDyiFt9MJvyxpUnd3bEfI99esM
YKH2JVSrt6q/VQBlD6Bkt/YMBIHRSzyY5ByeN69UjvxTEoxnRFQUIliRlgfpHV+I
Hr8ZsIOyT8PaQ8vT4bJns8A933VjQvX3r57mqiGz1yoW8WPkkXF8myrtF1Vqyf2p
E7oIvYk8cBddHVP7DYpWRipYkPg4zgi2+d9Mu8Nk2sAH6iayFxYPpkxWbFmPqEgU
yBjK66zCrHAgZ5AotyR0QNKjKJKV1bfUxXiiZ8h8RNBTEoxnRFQUIliRlgfpHV+I
yLdkXTDBy8zVmn/u7wrqvPxVe7oEZDSWD+/daEvtQPGEKm0qHVT/v4l7Ce+MhTug
qHKrmFqwU+R9vq/VBrUO9vZu/JXsX/fpQJhVwFhLACoawtVrndnNVlpZXc5teSwH
FBZz6k50XLLt3fqPn+9ej0BRQQBGo3FRGOSmtt0sMq4u4LaDKG4GqjcBIKy0Yr4g
XxYSvlFWhlieqpfDWHDRFkAeWD1n9JYlZ84jZSqRUTV0Afh6S76p4SdgnRIyoK9g
SjbJmMo5tC/7DaKETHBPxZHphzwkX3LoRzwXAftUD+2uhv9jGGO36H30aMtQxric
KklxYKTdUnzaPgnnRIPXKApK90qP0RN/c7+KTUtZswb+TQZUVmuTPFQTN5S3WXAH
x8dliRHnAQxn63xF/rUhzkmu31QKI1VkkC7mcbELs0DWOWHUuSjlFQjwgks8FSQa
1jlh1Lko5RUI8IJLPBUkGtY5YdS5KOUVCPCCSzwVJBrWOWHUuSjlFQjwgks8FSQa
1jlh1Lko5RUI8IJLPBUkGtY5YdS5KOUVCPCCSzwVJBpzvGAPOTaDkudLX8C9xcLs
nJLvY7wtrMnz9xgGfo2j4FGZ9zC0Ohz/ifpe1SbP1pc=
