# profile.d/xdg-environment.sh
U2FsdGVkX18TrSLuiOpjiM+OowyczGnr1EF4Mwj55qzYLbw8eyD97CIdH4ILH/Cd
RHkNL7GB2X+LnEuiHkyyMD/JfqSPZv4JUL2I4/SnRKoxJWR//RnDGszmGJ+Kcrga
AiYucVRq3xq2r0UZPszgDm7t/NjMHeCk4Eb56lXk/TZvIchCCWUPTdpAfhoA0+VY
5K6Qtq/OsB0gwsWElSFYctiCSbKNLQ0ocuzsGzQyddWfiiDxXLaxNVySpI4DTBTS
umdLVIJcJtC+qmlmMkjquU4HAcjUc0m3D/CiEPxeZyWeFqf5IYTB9I03GKZtDSoP
+9hNxWS5JUO6GSUpwJa7vyTLnDAIulrQqvHoXoBnCQepp2pbQm5d57I9h1pqZGlV
hyLHhRfPcF9InL26eJbwmKENkOpva0t5XH3FXQbizrCzsj3nPhj/q1GMozZRGLkI
Q/dDdYN8gYfELKRuMyLPl0O+4NFGC7rrT6kQoHCSZJSboG672ANpj1EiVUCj26yo
3ZJidd4UyNggAf+hdPWzall9SQiSgZIN3lQGZlE8qBLCy30DO3/1acr/WDS4ieVl
yaUTf68CGVyonbQx2iLZEkvzvL+qYL8EnmE9ADijRjwb4gnk9MZ5uRVpUfiTMgn3
di/05awJfdjUoiynmYWov7869YRnCdvd7850GL+KM5yCybgkzYtSt/ibmTYhQujo
I1jQMd2ITbdbq/dx+1t772s4+uXKPIlBdRw3yXfi4m6CwT/TcI9T424dWXm6kpOn
Lwq6sNehuGxpk25v2XNBRk9ML0K89s4uTliG40HsR2wKaLpMRR/w1U/BAq6DXmpl
iNijgs1hkJjDl2rq7hDx39yAGiKdL4fKdTf+fbsyA24kPIQvQQSPkIZB9X8SS/7z
xZGIt7lHJ8N5hUJAsBnBEbCKFas2ahRguTKF4IcdZjivXo22i2tSd+oK8SEItGca
jIJRTOa4/ILV4fgTuu9EydMbEG2JfVrOBONdUtRYEZtGwcmzREdzoh1YdwfoALZ2
FTZfCGDrbqvmAV7AMGbQsP6UpW6vD9J45I8cVYmv4JRtuvniwuLc4zSABz7unQa5
l0N3jOFMYtmQi5PuTniTcCnj+0LsabJkhf4RL49x+3adUd/PKDCSGs9lzaB1xk5v
hgMAaW5z9X+hfan9guGfmkfoaWxkD0PLU46KtOBoBh8kzT2qEyTXkrZPyoQyS9XD
R4cwYIDvdh1nymxovTdlTGlguTwZmnKu/5QKKyo3rye0Nd55IrgWQY+fY6ClY9In
cKGiaiGkIW0+0g6oBnOgubjHSPt5xnTW+Ql1ZTxkeRcRsFzBgotHevsmppQk6AoV
R4cwYIDvdh1nymxovTdlTH9L6DsM7drlXUu08MXDM2K8jwTSK8bBfFfop5VvrQAi
8kmjWPJbq5SGVW8FU7KQIenPWGVuO+UcT1d7py59r1lpDI//iwvYaGH4mjr3VbCs
VKQ4qZfuFU2QHhWLi8WSTrQ13nkiuBZBj59joKVj0idaSBUxCt3uXIgxhjbvCRRb
