# profile.d/xdg-environment.csh
U2FsdGVkX18TrSLuiOpjiPn6/0sn/gty/30C/HnpG99uVFg5a66jKMKsRd15DaHI
a5HvDZVhZFNhb7PBWrnWjg0WmyYyLDiVA+9hlGSe+zP0YOW4tttEN0m5pkWy2rRK
BXK2yw2vcq+1uEl5YXaS1TRN80jTM2y/z6+a5/BTFicANefWQqqdCEZ2/DDqifOt
MASeopg8VgC9Ss1ckqYW3YemVTXF4Mnf9+OrLcWEw4IogotcT6/mOT/qfLJa7Af2
JDZe6/QIEgAxFTrzMP4bgTG1axh0VErVAr6UNywnVNk8/iuK1YTRnRGdGTC3Yqrf
lqkJM0OgaLx0mYplfWV8yZgflFzyIjJU3qtUQZb0imZayOmORSa7C/4WgkyzF4k9
pSBX+SDXj+HwG1jtMNvmsN9s3NcJHEEwyoeiAlfnlrGmymiwMIdnvZI2+wHHdbQx
1knHEyyw12XWKGr5sda7+dJZcR6zQ1JoaeljI0DiJZfZsiESPw/K9QkH9bo1eFW+
HBThm/GXf3NPLgMHmCFsd25ObCFhvGfJwoERaA9zT5yqhjp74vVe0DCWixZ/Zz9h
z5q/V4GSZUboFvRa9hHs+BNTvkIOhezm+MlWxvEtmAN5rHWPBJVN8iCHhAyn7Bpe
MnX7p36XYD4ZLigPSCzwiKLDr/f31tjq8l4PRrCi15ANFpsmMiw4lQPvYZRknvsz
w/mjDJgP1cljwieQcAuIuPK0nsSiaw3NGjzjDMCgq2xpnVRouphyrN/m/Yb4eZU5
2yC8pXWrxDuiGza90f36BoFKeXxVqpZ5PFb0VE/8D4qMU6l4VGsEHQkYQ3pk0YiW
MU7/JfwTo/7AkrprVyTv5KFUaOkwDQoaPovLU4pRu7UM44EFwcHvPr7KbZfqz6pE
VXKNr6DkTn68ac3y7DjhpMpLscmdF/LlBXK/X12TwSd9tb/OXSeLMYVpA/xZSymf
FSm4OEo8acpBhdXoIDuv86xr8npecdxB/sf5KlysywocPY47qzcjhtA3lBNYNXnK
fOzchAsmDOwFxZBFDNuz3MiGnb2/+xWiZOBUJRLg0cS8gF6e89RMumr9686vWoeq
wMdOmNRQX+QYSdW6NJ/bzNgtlEZndBphYm+yPzbq+H/ql8FipXCFFKyJLPSH9j5Y
ZM/SumvNyJbaJzfRC2PbtNSEZ6QFCPMvEkN2xPX+cg4=
