# profile.d/zzz-groff.sh
U2FsdGVkX18TrSLuiOpjiHNX+2pLsSSgzcQaSoOwiBz+5ZgSwODQcTSeES9E67Tq
Wri9V2C+Ur8Pjp9fO7hNbvuPQi8k8xRPnDsY3sbdjQ52ZPlp+Le+RQ1wor+/w7Uv
sEPDL96yg1L1rIGcYapiZXVzX6+IDksLzzIVlINRKuo5jZufmyt2mOnZ0yBDz7Yl
w3zVggIrtfV1SeldV4MGzVdu84xFcmT6GUPgddEEV0YFW3DbMOYQXfegfMed/qqq
zucYBhHEf0C1Kuy9Khlmk7f+plLUb1ob+PtuA2l0coMa7xGRlEsjZ7hBuyx40yiB
30ZR6zmFFI0m9nIBgfnhTaU0kiTE8qFQbNIfJ0MjfmRaQXRBg9zAvk8AEiSG3eZl
