# profile.d/zzz-groff.csh
U2FsdGVkX18TrSLuiOpjiHNX+2pLsSSgzcQaSoOwiBz+5ZgSwODQcTSeES9E67Tq
tH2rywa+oIHd1vMDqJzpAPTs/bWw+Rgob/Pfzk+Mj1skRLT+jLMxQSJWZvKJ6jBO
VSLq3009Gcz1Ct+duX1wg5WnxlWtYsmElQeQRoyU2VJ+zNpb6r8lYYbLHaeM6zi/
jXFqj+ma+EholBbQ24kQxaDnpM9SlnUJsa2f1IHbSeZykna8WrrA+NgZ9T9vUeP4
NnAbt5hXV8rtWA7AKCq1a92aGIxwQBo4eI1rTq9YWhvSYXbiipKz7JyPAP1GA/a7
pEFVwHDXb7Eg0KubEZoQJe17nqYZqUDg1JLSqcEUrO3AbuPpp02S1/tlXBGboylO
KOZ5kpAAaltzCc2QwHJ5L7XL6LMzASlM0bJ18WBJC0v1MHhLvnTBJ3pzBegNLs2M
1xKj0kG32AsM6W6bkLO+xlcDZEUf3cRX3GEl4+cXRVS3PcX+w4MSck5TqeqWNLZr
eaKI/P1/HuELkKtD1XEBrw==
