# xinetd.d/cups-lpd
U2FsdGVkX18TrSLuiOpjiGeo5+2yftneoo89Yu4iyxqgF1Hm/HpKhPiEfN2d6Qx4
F9MS5Vg9qHhnezwmuA2WvQk1uriP2oMA4EVkIB7imE9HdYvAxizfDO6UeNtLIm6j
4i9KpEOjkHpYDq56oXbm3sLwOLB6Sihryb9U9k1lZ3+DE/RcoiRvUibe1DfYeLbd
gSHY1GgfTqRkJNTPN4Hy+tqL9uaiyQG6BPPvBM0kDeiLYkR+dCo2yjsNFD3l64Y1
miuS89iRl/2jF6vU4ljOYV3DKri/TYhRIjGJLjl7DjXveTnn1dhQmTadPksJsfsq
s8q5eoBuiJStBZxIGmztebR/jSa3UVF4Dph136AsCqZaQXRBg9zAvk8AEiSG3eZl
