# xinetd.d/cvs
U2FsdGVkX18TrSLuiOpjiBfknR2z8B3Q99G1NnbykU59enpUFABz0pyKNbbFsBro
fnZ6SyFvTh7tQdCGjrqFY3rEHv0fjdGSpXNCZLaV/F8ccCVIhfyn4UK4CgNH2af9
GUwhp0SFmDrfpLofnQ31mXmys59v7179skoJqtJ23bolZ9UZ4Qy1+SQDlSqVvvzB
fIzi1ZaA8fh+UVhaOoz9V4YOgltnUw14aWZlAfw08KWvnJYNOk+pq5NEYNIUxLvj
vTWhZgKRCU547d1JwneQcOEJ4SnPzHvjLhMOKtwBkGPAHT/1cKdIdROTX0G9Hmj5
qQ3i3at9NlzseIWSUYmZKKD1FWb5B3M7osn/DKskc9Kb0KoehD/hO5ciKrTEx5N7
Szj3yOpQU9ytmF05z55gVWJMLc37MNgCEjfEbeEJAqRYC4LehkfBB0Il/6zTG05J
+3hqE67n6BCpO7/Yl1Zc+qRn9m/BgO/MSvR87mb8BvS/rHKDYmoqasx3jnwxSwd1
tXT05MnwRJ/yfvQnAi0s9hJCLSiV26Yf43OlkeuoIkB7lOLl7kDMi6j9HkVww/y/
