# xinetd.d/discard-udp
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvurS/waho0VIGbiwi3qtgHm
bWmxHQbQM6s3E56VZf2EdGPilbCbttgyS2Q0r67YLLV2NF9aLrLAXmL57NMnCV8w
11fa1lTkKBTJOvTsEXKpMKsJCj1zQypay+by5Qxx4xJa5WbGgQNB5CeK5kSpcibG
p7p/b+uN3L1ONS9ma27YOWtZZmF+RXmtMyGuSXwp15Iy/sr2/GP6pw2/t2dELzLs
U+tafZOg3RJmgZNhc3tCThbuS8PYcJGfRrsvy2XXh/mNJZChq4oz9/uSjTdNmrFb
zyP5+cjukE2fp/ixk9my1ASanihEPujwU2KAiLzk3pA4oANf3IyGhxL5AfJ4G+0r
OKADX9yMhocS+QHyeBvtKzigA1/cjIaHEvkB8ngb7Ss4oANf3IyGhxL5AfJ4G+0r
Np1S0oCPlqjLudlp6ogtmg==
