# xinetd.d/echo-udp
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvsSeuv4Jmby+UsK92+cOXp+
hcSE/w1sUp8ie7nWGLEq7C41ADdYcVqDLWgamIODToWpYeXrbilWyL+Yj5Zu0YmC
nUOME+3qKr3Zua8M0fgr8k9VxQljGAGW5tkKIex1MH1bFQs4imUN7jwL9Qmf0tTw
lqckxu1y1eXlQre0hJhsvfWGr7RqJb8nEMdvFuIDVQhllGScq/g+dLzZzeEeSzn7
S44K94obDK+sbMsabcEzllCbskx5/dKfAF1iJVT7Fo4E4+ZRSpv+oa1AG3ktvGIQ
yZcLlh0jeUk9s+N3C2e6rlE6btdXv/IjGxEcxj/RlVw4oANf3IyGhxL5AfJ4G+0r
OKADX9yMhocS+QHyeBvtKzigA1/cjIaHEvkB8ngb7Ss4oANf3IyGhxL5AfJ4G+0r
YT0kmrwvDaNYHe5+V9CspQ==
