# xinetd.d/rsync
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvs/XQQJqjdcPjtfN2cFBTxM
zz3sfMCOPxu9ZsH6QkU/Audt0nMMqpaxN0cZNMLwS9zneEk193CBuemIfpi3iTEh
490WAgCuYaWyjjw+KXV4hg7lnLTI6TAMvXPvwU0joNWTjvL3BBh4aknmJT1ACAgN
hbUcLKAoB6BWB8OYqilDSkBrHucgyIOrQXIZzrJ366Cs+iQZPsvJFupEacHdUuI5
HoDQO4302qc272rGf5uIXdopBgzMBZ6Uh4m9mTyXzU0=
