# xinetd.d/svnserve
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvtL5W6MBpji+rahNkatSjr/
LU1/nYek5xrEgHWUxgJCev7Fo7cUmryCINYt8jPZYS9R8Lmhcv5lbrgbCxYsZpxu
/y8BlzXDHxGE+SIGku/iwPnmbJ4KjBrOY0p3oExO0Kryo/9OAVbPRev2mmQwwem0
+rYw9SS5Yz3iwhYWCytCtABGDSmqhejRfSOSjJ5xbU6bHu9uk7gLEyTp62zuE29w
9HxsA7akCvmGXt7CX2MhHdIhVs/R2oOG2lNE3wyWsQwfCbsGIXec7OnpawP3WTsZ
zWeI4XgnmLLrcPlg1YQw20uN0pvb04z/1D9I9qu/TWfusDe2PUkS61+aKNIeZ0JE
NcxSr7JXL7t7Tr1Tl3bAVIhViYGTgwV+R4eMIvkvYWYEUacPxQdSjLeINVN/zEpw
gsHH0pZaC8i9LuOrHkL1vB51JzOp6hDNfmn1V5OgCH4=
