# xinetd.d/tftp
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvvmMIX6a3LgOgpy6221Qgve
uBM1GLR73UzMZwQx+4+s8qCF0227Cl7uT9u+qQT52kjvS43pyuWE6Anw3zpTFHXh
RoW2Usj2bCbR5wMNi+T6GjJ+k051Pq37aYWXoM2kA0rybSWjfsCITWlmU7UpS14P
xHgDYOzaulWL3/F8/pXw+mzTJAZg1PAIpGP1/h+JqEKgdduQai+GkQuNxGDe2/x/
sizmTKWSJ3nEYtmzmsqWmRO/cEQBF3IG7M457y84uxmZbbP/F9002D1DBSq1YbNd
vLNjGMJNLEPkJv5NMRL2gxK5CyNQzmYP14OXPPGhFbYoCIlioDLOkMskl9GmIOaW
SokeAGiyIwx/sfU1ZzYP6Ps7OjEYkCna5hTQBc4dTSm/t59DuvF2MRWCjGsFv963
9l3KTMkwrF5Fh6Ou+T6PkXHdEUOgMZJrMlhrxFig2tfjVm9tftojQEggklbhnHi2
gFXDAK0BcN8LEFcbe1voK60lOO6PeZP0kRnterRB8T6MqNpluqjGKfGn3U8NCp1c
iLIWw+6NprA0jUY7A8kAsleV+vzmf+vVYb9Di30ufygdwCeWo2UoWWo0wGAnU1mQ
9RQTGDKAoIBimD8HQFt/EPc/pGWwJPfbSz5ym6aAabn7KKYK++W0Eh48ziPckvE9
eefn8qYn7QZke3JIdiWDBIsmPLMBVz/XjBxCNqq/cgSf2VmBd1Gy0MomPlxPJ2Ks
zx8Zy4+Ix/xkPb0/Fs/Q0qoi+IUy6XfCRhNWkxxCXigutxEH/4AI7alle4pzO2zT
