# xinetd.d/time
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvsSeuv4Jmby+UsK92+cOXp+
FIgEgIRqERo2eAIVS+oNPpGIB7az6cTq2gLsGAKFYww2OMx3CWS8NqULxUF/e8VD
C5jLPKZRL2fNsEuTOjFUOwVzp1x+LMH84CNZrb7JXhtcQg+OgSPt6W2In1jB6BWD
Lse/dUsi5Eo903ZO4orwQQJ0UWXPWv+N13wkdkzlQfF2Uqf108JgZwN34C32mY+J
SPMW0h+X7y755loR6giUwbuuuZ92i9L/1cUMeKzd41u5TtafsZb9jTqJnYpRfzpc
c+l94fclbsJHOfbgQ1inA8z2cQLbeWJl3YoZA/oe5b7FW4RJ7Z0S0PVd/yXNmy7Y
OKADX9yMhocS+QHyeBvtKzigA1/cjIaHEvkB8ngb7Ss4oANf3IyGhxL5AfJ4G+0r
OKADX9yMhocS+QHyeBvtK6QrOZCPFHURUKWGeiLJ57I=
