# xinetd.d/time-udp
U2FsdGVkX18TrSLuiOpjiNmqpyrR9eD4vIMKoMIjTvsSeuv4Jmby+UsK92+cOXp+
FIgEgIRqERo2eAIVS+oNPpGIB7az6cTq2gLsGAKFYww0vwadBnvATKpDBthLRj18
LmyMauX35zSiqZRouUkrstBbo6rrZ0N4C1z8grSvylN8kuZoHiTvMwn5ARYdxMbJ
/4Ryvy0k45OeTi2JKLHnR2tZZmF+RXmtMyGuSXwp15Iy/sr2/GP6pw2/t2dELzLs
U+tafZOg3RJmgZNhc3tCThbuS8PYcJGfRrsvy2XXh/mNJZChq4oz9/uSjTdNmrFb
JylpmWy1RdvR1rNXxDpMCoQJGgUOAf0F7WnUL54PQTw4oANf3IyGhxL5AfJ4G+0r
OKADX9yMhocS+QHyeBvtKzigA1/cjIaHEvkB8ngb7Ss4oANf3IyGhxL5AfJ4G+0r
hJxHoTbDLZOBRnDU421K6A==
